# fix/disable-nic-aspm.timer
[Unit]
Description=Re-assert NIC ASPM disable periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
