# fix/nic-detach-watchdog.service
[Unit]
Description=Power-cycle when igc NIC falls off the PCIe bus

[Service]
ExecStart=/usr/local/sbin/nic-detach-watchdog.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
